# crontab
* * * * * python3 manage.py wallet-balance-snapshot
