# crontab
* * * * * python3 /app/manage.py wallet-balance-snapshot >> /test_file 2>&1 > /proc/1/fd/1 2> /proc/1/fd/2
