# crontab
* * * * * python3 manage.py wallet-snapshot-data
